# car-pi-control.service
[Unit]
Description=Pi Control

[Service]
Type=simple
Restart=always
RestartSec=1
User=root
ExecStart=/home/obito1903/car-pi-control/car-pi-control -timer=10m

[Install]
WantedBy=multi-user.target
